# solar_initial_bash_service(deprecated)/version03_2021/ladex10.service
[Unit]
Description=Solar Script

[Service]
Type=simple
ExecStart=/usr/bin/ladex10.sh
 
[Install]
WantedBy=multi-user.target

